# scripts/fancontrol/fancontrol.service
[Unit]
Description=Fancontrol

[Service]
ExecStart=/usr/bin/fancontrol

[Install]
WantedBy=multi-user.target